# deploy/systemd/manavault-remote.service.example
[Unit]
Description=ManaVault authenticated remote access
After=network-online.target manavault.service
Wants=network-online.target
Requires=manavault.service

[Service]
Type=simple
WorkingDirectory=/opt/manavault
ExecStart=/opt/manavault/.venv/bin/python -m uvicorn backend.remote:app --host 127.0.0.1 --port 8002
Restart=on-failure
RestartSec=5
User=manavault
Group=manavault
Environment=PYTHONUNBUFFERED=1
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=full
UMask=0077

[Install]
WantedBy=multi-user.target
